# debian/heartbeat.files
etc/ha.d/README.config
etc/ha.d/rc.d/ifstat
etc/ha.d/rc.d/ip-request
etc/ha.d/rc.d/ip-request-resp
etc/ha.d/rc.d/ns_st
etc/ha.d/rc.d/shutdone
etc/ha.d/rc.d/status
etc/ha.d/resource.d/AudibleAlarm
etc/ha.d/resource.d/Filesystem
etc/ha.d/resource.d/IPaddr
etc/ha.d/resource.d/MailTo
etc/ha.d/resource.d/Raid1
etc/ha.d/shellfuncs
usr/share/doc/heartbeat/AUTHORS
usr/share/doc/heartbeat/ChangeLog
usr/share/doc/heartbeat/GettingStarted.txt
usr/share/doc/heartbeat/GettingStarted.html
usr/share/doc/heartbeat/HardwareGuide.txt
usr/share/doc/heartbeat/HardwareGuide.html
usr/share/doc/heartbeat/README
usr/share/doc/heartbeat/Requirements.txt
usr/share/doc/heartbeat/Requirements.html
usr/share/doc/heartbeat/faqntips.html
usr/share/doc/heartbeat/faqntips.txt
usr/share/doc/heartbeat/rsync.html
usr/share/doc/heartbeat/rsync.txt
usr/share/doc/heartbeat/authkeys
usr/share/doc/heartbeat/ha.cf
usr/share/doc/heartbeat/haresources
usr/lib/heartbeat/ResourceManager
usr/lib/heartbeat/api_test
usr/lib/heartbeat/findif
usr/lib/heartbeat/heartbeat
usr/lib/heartbeat/mach_down
usr/lib/heartbeat/plugins/HBauth/crc.a
usr/lib/heartbeat/plugins/HBauth/crc.la
usr/lib/heartbeat/plugins/HBauth/crc.so
usr/lib/heartbeat/plugins/HBauth/md5.a
usr/lib/heartbeat/plugins/HBauth/md5.la
usr/lib/heartbeat/plugins/HBauth/md5.so
usr/lib/heartbeat/plugins/HBauth/sha1.a
usr/lib/heartbeat/plugins/HBauth/sha1.la
usr/lib/heartbeat/plugins/HBauth/sha1.so
usr/lib/heartbeat/plugins/HBcomm/bcast.a
usr/lib/heartbeat/plugins/HBcomm/bcast.la
usr/lib/heartbeat/plugins/HBcomm/bcast.so
usr/lib/heartbeat/plugins/HBcomm/mcast.a
usr/lib/heartbeat/plugins/HBcomm/mcast.la
usr/lib/heartbeat/plugins/HBcomm/mcast.so
usr/lib/heartbeat/plugins/InterfaceMgr/generic.a
usr/lib/heartbeat/plugins/InterfaceMgr/generic.la
usr/lib/heartbeat/plugins/InterfaceMgr/generic.so
usr/lib/heartbeat/plugins/test/test.a
usr/lib/heartbeat/plugins/test/test.la
usr/lib/heartbeat/plugins/test/test.so
usr/lib/heartbeat/req_resource
usr/lib/heartbeat/send_arp
usr/lib/libhbclient.a
usr/lib/libhbclient.la
usr/lib/libhbclient.so.0
usr/lib/libhbclient.so.0.0.0
usr/lib/libpils.a
usr/lib/libpils.la
usr/lib/libpils.so.0
usr/lib/libpils.so.0.0.0
usr/share/man/man8/heartbeat.8.gz
var/lib/heartbeat/api
var/lib/heartbeat/casual
var/lib/heartbeat/fifo
var/lib/heartbeat/register
etc/init.d/heartbeat
etc/logrotate.d/heartbeat
