# debian/heartbeat.conffiles
etc/init.d/heartbeat
etc/logrotate.d/heartbeat
