# debian/heartbeat-ldirectord.conffiles
/etc/logrotate.d/ldirectord
/etc/init.d/heartbeat-ldirectord
