# debian/heartbeat-ldirectord.files
/etc/ha.d/resource.d/ldirectord
/etc/init.d/ldirectord
usr/share/doc/heartbeat-ldirectord/ldirectord.cf
etc/logrotate.d/ldirectord
usr/share/man/man8/ldirectord.8.gz
usr/sbin/ldirectord
